# workshops/lab-lookup-consumption/workshop/supervisor/web-server.conf
[program:web-server]
process_name=web-server
command=/home/eduk8s/exercises/web-server.sh
stdout_logfile=/proc/1/fd/1
stdout_logfile_maxbytes=0
redirect_stderr=true
autorestart=true
